# workshop/SANDBOXING/systemd.units/chroot.service
[Service]
RootDirectory=/opt/debian
MountAPIVFS=true
