# workshop/SANDBOXING/systemd.units/dynamic_user.service
[Service]
DynamicUser=yes
ExecStart=sleep infinity
