# workshop/SANDBOXING/systemd.units/dynamic_user.service
[Service]
ExecStart=sleep infinity
DynamicUser=yes
