# workshop/SANDBOXING/systemd.units/mount_ns.service
[Service]
ProtectHome=yes
ProtectSystem=full
PrivateTmp=true
BindPaths=/usr/lib/systemd
ExecStart=sleep infinity
