# workshop/SERVICES/systemd.units/life.service
[Unit]
Requires=water.service
After=water.service

[Service]
ExecStart=echo wow
