# workshop/JOURNAL_101/systemd.units/journal_hello.service
[Service]
Type=oneshot
StandardOutput=journal
ExecStart=echo "Hello workshop from service unit!"