# workshop/SYSTEMD_101/systemd.units/bone.service
[Service]
Type=notify
ExecStart=/bin/false
