# workshop/SYSTEMD_101/systemd.units/dinosaur.service
[Unit]
Before=human.service

[Service]
ExecStart=echo rawr
RemainAfterExit=yes
