# workshop/SYSTEMD_101/systemd.units/dog.service
[Unit]
Wants=bone.service

[Service]
Type=oneshot
ExecStart=echo woof
