# workshop/SYSTEMD_101/systemd.units/human.service
[Service]
ExecStart=echo oogabooga
RemainAfterExit=yes
