# workshop/SYSTEMD_101/systemd.units/noodle.service
[Service]
ExecStart=/bin/echo "very noodley end"

[Install]
WantedBy=pasta.target
