# workshop/SYSTEMD_101/systemd.units/pasta.target
[Unit]
Description=My Very Pasta Target
Requires=multi-user.target
Conflicts=rescue.service rescue.target
After=multi-user.target rescue.service rescue.target
AllowIsolate=yes
