# conf/profile.d/workshop2path.sh
if [[ -e "/usr/local/src/workshop" ]]; then
    export PATH="/usr/local/src/workshop/bin:\${PATH}"
fi